# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/root/logs/supervisord.log
childlogdir=/root/logs
logfile_maxbytes=2MB
logfile_backups=10
loglevel=info
directory=/root

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:rethinkdb]
autostart=true
autorestart=true
command=rethinkdb
stdout_logfile=/root/logs/%(program_name)s.log
stderr_logfile=/root/logs/%(program_name)s.log

[program:turtl]
autostart=true
autorestart=true	
command=/usr/bin/bash /turtl/run.sh
stdout_logfile=/root/logs/%(program_name)s.log
stderr_logfile=/root/logs/%(program_name)s.log
